# services/connect/run
#!/bin/bash

cd /app/
exec crossbar start